# config/serve_logs.conf
[include]
files=base.conf

[program:serve_logs]
command=airflow serve_logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
